# autostart/Eww.desktop
[Desktop Entry]
Type=Application
Name=Eww
Exec=sh -c "eww open omni-tray"
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
